# 0.60/helium.service
#helium_gateay.service
[Unit]
Description= helium gateway

[Service]
type=forking
ExecStart=/root/update/helium_gateway server