# 0.60/helium.service
#helium.service
[Unit]
Description= helium gateway

[Service]
type=forking
ExecStart=/etc/helium_gateway/helium_gateway server